# root/etc/services.d/ombi/run
#!/usr/bin/with-contenv bash

cd /app/ombi || exit

exec \
	s6-setuidgid abc dotnet /app/ombi/Ombi.dll
