# root/etc/services.d/ombi/run
#!/usr/bin/with-contenv bash

cd /opt/ombi || exit

exec \
	s6-setuidgid abc /opt/ombi/Ombi --storage "/config" --host http://*:3579 --baseurl "${BASE_URL:-/}"
