# root/etc/services.d/ombi/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc mono /opt/ombi/PlexRequests.exe
